# scripts/deploy/piddy.service
[Unit]
Description=Piddy AI Backend Developer - Slack Integration
After=network.target postgresql.service

[Service]
Type=simple
User=%i
WorkingDirectory=/workspaces/Piddy
Environment="PATH=/workspaces/Piddy/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin"
ExecStart=/workspaces/Piddy/venv/bin/python -m src.service.background_runner start
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/workspaces/Piddy

# Resource limits
CPUQuota=80%
MemoryLimit=2G

[Install]
WantedBy=multi-user.target
